# config/live-build/config/includes.chroot/etc/systemd/system/getty@tty1.service.d/hide.conf
[Service]
TTYVTDisallocate=no
StandardInput=tty
StandardOutput=null
StandardError=null
