# config/live-build/config/includes.chroot/etc/udev/rules.d/99-horizon-power.rules
ACTION=="add|change", SUBSYSTEM=="power_supply", KERNEL=="BAT*", RUN+="/bin/chmod 0666 %S%p/charge_control_start_threshold %S%p/charge_control_end_threshold %S%p/charge_start_threshold %S%p/charge_stop_threshold"
